# templates/default/services/source-mirror.conf.erb
#!upstart
description "Update source mirror"

start on started source-mirror-shim
stop on stopping source-mirror-shim

task
respawn

setuid <%= node.source_mirror.user %>

env USER=<%= node.source_mirror.user %>
env HOME=/home/<%= node.source_mirror.user %>

exec <%= node.source_mirror.data_dir %>/source-mirror.sh
